# deb/deb-src/sysroot/etc/init/bigboards-fabric.conf
description "BigBoards Fabric"
author      "bigboards - hello@bigboards.io"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on filesystem and started networking
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/usr/share/bigboards/fabric"

    exec sudo -u bb /usr/bin/nodejs /usr/share/bigboards/fabric/server/server.js >> /var/log/bigboards-fabric.log 2>&1
end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
   echo "BigBoards fabric has started on port 7000"
end script