# deb/deb-src/sysroot/etc/init/bigboards-fabric.conf
description "BigBoards Fabric"
author      "bigboards - hello@bigboards.io"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on filesystem and started networking
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/usr/share/bigboards-fabric"
    export BB_FABRIC_CONFIG_FILE="/etc/bigboards/fabric.yaml"

    chdir /usr/share/bigboards-fabric
    exec sudo -u bbfabric /usr/bin/nodejs code/server.js >> /var/log/bigboards/bigboards-fabric.log 2>&1
end script